# deploy/offsite/pbs-offsite-sync.timer
# Runs the off-site sync nightly at 04:30 — AFTER the local (03:00) and PBS
# (02:30) backup jobs finish, so it ships a fresh, settled datastore.
[Unit]
Description=Nightly off-site PBS datastore sync

[Timer]
OnCalendar=*-*-* 04:30
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
