# deploy/offsite/pbs-offsite-sync.service
# Off-site PBS datastore sync (one-shot, triggered by the .timer).
# Install on the PBS host (node 2):
#   cp deploy/offsite/pbs-offsite-sync.sh /usr/local/bin/ && chmod +x /usr/local/bin/pbs-offsite-sync.sh
#   cp deploy/offsite/pbs-offsite-sync.{service,timer} /etc/systemd/system/
#   systemctl daemon-reload && systemctl enable --now pbs-offsite-sync.timer
[Unit]
Description=Encrypted off-site sync of the PBS datastore to cloud (rclone)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/pbs-offsite.env
ExecStart=/usr/local/bin/pbs-offsite-sync.sh
Nice=10
IOSchedulingClass=idle
# Don't hammer the box if it loops; the timer will retry next window.
